# scripts/systemd/krx-price.service
[Unit]
Description=krx-price Data Collection Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory={{PROJECT_ROOT}}
ExecStart=/usr/bin/docker compose run --rm app
User={{USER}}

[Install]
WantedBy=multi-user.target
