# scripts/systemd/marketdata-fetch.timer
[Unit]
Description=MarketData Collection Timer

[Timer]
OnCalendar=Mon..Fri *-*-* 17:00:00
Persistent=true

[Install]
WantedBy=timers.target
